# crontab/crontab
*/5 * * * * /retention.sh >> /var/log/cron.log 2>&1
